# systemd/user/pipewire.service.d/override.conf
[Service]
ExecStartPre=/bin/rm -f %t/pipewire-0.lock
